# util/cron/nightly-gasnet.fast.cron
#!/usr/bin/env tcsh
echo `date` "" Starting $0 on `hostname`.

#
# Test a GASNet compile using the fast segment
#
source ~/.chapelcshrc-cron-gasnet
setenv CHPL_GASNET_CFG_OPTIONS --disable-ibv
setenv CHPL_GASNET_SEGMENT fast
set LOGFILE = $HOME/tmp/nightly-gasnet-cron-$CHPL_GASNET_SEGMENT.out
echo `date` "" "Testing GASNet using segment $CHPL_GASNET_SEGMENT.  Log file: $LOGFILE"
$CHPL_HOME/util/cron/nightly -cron >&! $LOGFILE

echo `date` "" Finished $0.
