# util/cron/nightly-no-local.cron
#!/usr/bin/env bash
#
# --no-local testing 
#
# If -compperformance is passed as $1, will run nightly with 
# -compperformance (--no-local).

CWD=$(cd $(dirname $0) ; pwd)

source $CWD/common.bash

if [ "${1}" == "-compperformance" ] ; then
    compperf_opt="-compperformance (--no-local)"
else
    compperf_opt=""
fi

log_info "Testing --no-local."

$CWD/nightly -cron -no-local $compperf_opt

log_info "Finished ${0}."
