# util/cron/nightly-opts-examples.cron
#!/usr/bin/env bash

rm -f /tmp/nightly.opts.out
export CHPL_HOME_REPOSITORY=https://chapel.svn.sourceforge.net/svnroot/chapel/trunk && /users/chapelu/chapel/util/cron/start_opt_test -cron >& /tmp/nightly.opts.out
