# util/cron/nightly-no-local.cron
#!/bin/tcsh
echo `date` "" Starting $0 on `hostname`.

set LOGFILE = $HOME/tmp/nightly-nolocal.out
echo `date` "" "Testing --no-local.  Log file: $LOGFILE"
source ~/.chapelcshrc-cron && ~/chapel/util/cron/nightly -cron -no-local >&! $LOGFILE

echo `date` "" Finished $0.
