# util/cron/nightly-opts-all.cron
#!/usr/bin/env bash

rm -f /tmp/nightly.opts.long.out
export CHPL_HOME_REPOSITORY=svn://svn.code.sf.net/p/chapel/code/trunk && /users/chapelu/chapel/util/cron/start_opt_test -cron -all >& /tmp/nightly.opts.long.out
