# util/cron/nightly-gasnet.fast.cron
#!/usr/bin/env bash

source $(cd $(dirname $0) ; pwd)/common-gasnet.bash

#
# Test a GASNet compile using the fast segment
#
export CHPL_GASNET_SEGMENT=fast
export LOGFILE=~/tmp/nightly-gasnet-cron-$CHPL_GASNET_SEGMENT.out

log_info "Testing GASNet using segment $CHPL_GASNET_SEGMENT.  Log file: $LOGFILE"
with_logging $LOGFILE $CHPL_HOME/util/cron/nightly -cron

log_info "Finished ${0}."
