# util/cron/nightly-mpi.cron
#!/usr/bin/tcsh

source ~/.chapelcshrc-cron

setenv CHPL_COMM mpi

~/chapel/util/cron/nightly -examples -cron
