# util/cron/nightly-gasnet.fast.cron
#!/usr/bin/env tcsh

#
# Test a GASNet compile using the fast segment
#
source ~/.chapelcshrc-cron-gasnet
setenv CHPL_GASNET_CFG_OPTIONS --disable-ibv
setenv CHPL_GASNET_SEGMENT fast
$CHPL_HOME/util/cron/nightly -cron &! ~/tmp/nightly-gasnet-fast.out

