# test/Bin/nightly-mta.cron
#!/usr/bin/tcsh

source ~/.chapelcshrc-mta-cron
$CHPL_HOME/test/Bin/nightly -cron -examples

# MTA compile compiler
setenv CHPL_HOST_PLATFORM $CHPL_TARGET_PLATFORM
unsetenv CHPL_TARGET_PLATFORM
$CHPL_HOME/test/Bin/nightly -cron -notest

source ~/.chapelcshrc-xmt-cron

# XMT cross-compile test
$CHPL_HOME/test/Bin/nightly -cron -examples -componly

# 
setenv CHPL_HOST_PLATFORM $CHPL_TARGET_PLATFORM
unsetenv CHPL_TARGET_PLATFORM
$CHPL_HOME/test/Bin/nightly -cron -notest

