# test/Bin/nightly.debug.cron
#!/bin/tcsh
source ~/.chapelcshrc-cron
$CHPLHOME/test/Bin/nightly -debug
