# test/Bin/nightly.novgr.examples.cron
#!/bin/tcsh
source ~/.chapelcshrc-cron
$CHPL_HOME/test/Bin/nightly -cron -examples
