# test/Bin/nightly-pvm.cron
#!/usr/bin/tcsh

#
# Test a PVM compile
#
source ~/.chapelcshrc-cron-pvm
$CHPL_HOME/test/Bin/nightly -cron -examples
