# test/Bin/nightly-gasnet.cron
#!/usr/bin/tcsh

#
# Test a GASNet compile using the default segment (everything for linux)
#
source ~/.chapelcshrc-cron-gasnet
setenv CHPL_GASNET_SEGMENT everything
$CHPL_HOME/test/Bin/nightly -cron

#
# Then test one using the other segment (e.g., fast)
#
setenv CHPL_GASNET_SEGMENT fast
$CHPL_HOME/test/Bin/nightly -cron
